# Fluid/agent/reolicator_svc.conf
description "Fluid Lazy Replicator"
author "xzzhu@mail.ustc.edu.cn"

start on runlevel [2345]
stop on runlevel [0156]

script
CMD
end script

